# files/etc/hotplug.d/dhcp/99-update-and-notify-dns
#!/bin/sh

# esporta tutte le variabili di hotplug
export $(env | awk -F= '{print $1}')

# esegui il vero script in background
/root/update-dns.sh &
exit 0
